# config/fan-control.service
[Unit]
Description= Enable fan control

[Service]
ExecStart=python /home/scanner/fan_control.py

[Install]
WantedBy=multi-user.target